# recipes-metrological/webbridge/files/webbridge.service.in
[Unit]
Description=Webbridge
Wants=network-online.target
After=network-online.target @EXTRA_AFTER@

[Service]
PIDFile=/var/run/webbridge.pid
EnvironmentFile=-/etc/webbridge/webbridge.env
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run"
ExecStart=/usr/bin/webbridge
ExecStop=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
