# recipes-wpe/wpeframework/files/wpeframework.service.in
[Unit]
Description=wpeframework
Wants=multi-user.target
After=multi-user.target @EXTRA_AFTER@

[Service]
PIDFile=/var/run/WPEFramework.pid
EnvironmentFile=-/etc/wpeframework/WPEFramework.env
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run"
ExecStart=-/usr/bin/WPEFramework -b
ExecStop=/bin/kill $MAINPID
ExecPre="mkdir -p /home/root/Netflix/dpi; ln -sfn /etc/playready /home/root/Netflix/dpi/playready"

[Install]
WantedBy=multi-user.target
